# docker-files/incsac/config/etc/services.d/swoft/run
#!/bin/sh

cd /var/www/html

# php bin/swoft http:start

/opt/swoftcli run -c http:start -b bin/swoft /var/www/html